# packaging/etc/profile.d/latx.sh
#
# Loads LATX (Loongson Architecture Translator for x86) global options.
#
# Edit /etc/latx.conf to change LATX's default behaviors.
#
source /etc/latx.conf
